# config/etc/xdg/autostart/display.desktop
[Desktop Entry]
Name=jukebox pi
Exec=/home/pi/kiosk